# scripts/pp-main.conf
# pp-main Start the main Photopops button loop

description     "Photopops"

start on startup

task
exec sudo -u photopops /usr/bin/python /opt/photopops/pp.py
