# scripts/pp-node.conf
# pp-node Start node.js for photopops

description     "node.js for Photopops"

start on startup

task
exec sudo -u photopops /home/photopops/local/bin/node /opt/photopops/bin/popnode.js
