# scripts/pp-worker.conf
# pp-worker Start the main Photopops button loop

description     "Photopops RabbitMQ worker process"

start on startup

task
exec sudo -u photopops /usr/bin/python /opt/photopops/pp-worker.py
